# Cyclone/update-motd.d/00-custom
#!/bin/bash
echo "Welcome to Cyclone OS"
echo
echo "For support, contact support@iagon.com
echo
echo
